# monitoring/PerlModules/NP/Scheduler/init_script
#!/bin/sh
#
# Startup script for the NOCpulse Scheduling Kernel
#
# chkconfig: 345 95 15
# description: NOCpulse scheduler
# processname: kernel.pl


# Source function library.
. /etc/rc.d/init.d/functions

daemonUser=nocpulse
bin=/opt/home/nocpulse/bin
var=/opt/home/nocpulse/var
heartbeatFile=${var}/kernel.log
heartbeatFreq=300
heartbeatCheck=600
logLevel=${2:-1}

# Make sure we're in a directory readable by $daemonUser
cd /tmp

# See how we were called.
case "$1" in
  start)
	echo -n "Starting kernel.pl: "
	${bin}/gogo.pl --hbfile $heartbeatFile --hbfreq $heartbeatFreq --hbcheck $heartbeatCheck --user $daemonUser -- ${bin}/kernel.pl --loglevel=$logLevel >>${var}/kernel-error.log 2>&1 &
	echo_success
	echo
	;;
  stop)
  	/opt/home/nocpulse/bin/gogo.pl --kill kernel.pl
	exit $?
	;;
  status)
  	/opt/home/nocpulse/bin/gogo.pl --check kernel.pl
	exit $?
	;;
  restart)
	$0 stop
	$0 start
	exit $?
	;;
  *)
	echo "Usage: $0 {start|stop|restart|status}"
	exit 1
esac

exit 0
