# client/tools/rhn-virtualization/scripts/rhn-virtualization.cron
0-59/2 * * * * root @PYTHON@ < @PYTHONPATH@/virtualization/poller.py
