# client/tools/rhn-virtualization/scripts/rhn-virtualization.cron
0-59/2 * * * * root python < /usr/share/rhn/virtualization/poller.py
